# config.betting.cron
# 10 rerun cron
*/10 * * * * /home/g0d/Workspace/Projects/betting/cron.betting.play-gbets.sh >> /home/g0d/Workspace/Projects/betting/cron.betting.play-gbets.$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1