# crons/config.betting.cron
# 10 rerun cron
*/15 * * * * /bin/bash -c 'cd /home/g0d/Workspace/Projects/betting && python3 run_gbets_bet_split.py && python3 run_gbets_bet_combo.py' >> /home/g0d/Workspace/Projects/betting/cron.betting.play-gbets.$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1